# common-files/usr/lib/systemd/system/auroramax-firstboot.service
# AuroraMax First Boot Service
# /usr/lib/systemd/system/auroramax-firstboot.service

[Unit]
Description=AuroraMax GameHack First Boot Setup
ConditionPathExists=!/var/lib/auroramax/.firstboot-done
After=multi-user.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/auroramax-firstboot.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
